# deploy/supervisor_gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn --workers 3 --bind 127.0.0.1:9998 office.wsgi:application
directory=/opt/office
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
